# supervisor/supervisord_couchpotato.conf
[program:couchpotato]
directory=/opt/couchpotato
command=python CouchPotato.py
user = root
autostart=true
autorestart=true